# etc/systemd/system/wetterradar-windfetch.service
[Unit]
Description=Wetterradar Wind Fetcher (Europe only)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=www-data
WorkingDirectory=/var/www/wetterradar
ExecStart=/usr/bin/node /var/www/wetterradar/wind-fetcher.js --region europe --out /var/www/wetterradar/wind/current.json --ttl-hours 6

[Install]
WantedBy=multi-user.target
