# etc/systemd/system/wetterradar-windfetch.service
[Unit]
Description=Wetterradar Wind Fetcher (regional + global fields)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=www-data
WorkingDirectory=/var/www/wetterradar
ExecStart=/usr/bin/node /var/www/wetterradar/wind-fetcher.js --region germany --out /var/www/wetterradar/data/wind-germany.json --ttl-hours 6
ExecStart=/usr/bin/node /var/www/wetterradar/wind-fetcher.js --region europe --out /var/www/wetterradar/data/wind-europe.json --ttl-hours 6
ExecStart=/usr/bin/node /var/www/wetterradar/wind-fetcher.js --region world --out /var/www/wetterradar/data/wind-global.json --ttl-hours 6

[Install]
WantedBy=multi-user.target
