# etc/systemd/system/wetterradar-windfetch.service
[Unit]
Description=Wetterradar Wind Fetcher (Europe only)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=www-data
WorkingDirectory=/var/www/wetterradar
ExecStart=/usr/bin/node /var/www/wetterradar/wind-fetcher.js
ExecStartPost=/bin/cp -f /var/www/wetterradar/wind/current.json /var/www/wetterradar/wind/fallback.json
ExecStartPost=/bin/chmod 664 /var/www/wetterradar/wind/current.json /var/www/wetterradar/wind/fallback.json

[Install]
WantedBy=multi-user.target
